# deploy/north-coordinator@.service
# systemd template unit — one coordinator per tenant (instance-per-tenant).
#   sudo cp deploy/north-coordinator@.service /etc/systemd/system/
#   # per-tenant env (unique FRAM_PORT + that tenant's log):
#   sudo install -d /etc/north
#   printf 'FRAM_PORT=7801\nFRAM_LOG=/var/lib/north/tenants/acme/facts.log\n' \
#     | sudo tee /etc/north/acme.env
#   sudo systemctl enable --now north-coordinator@acme
[Unit]
Description=North coordinator (tenant %i)
After=network.target

[Service]
Type=simple
User=north
Group=north
Environment=FRAM_HOME=/opt/fram
EnvironmentFile=/etc/north/%i.env
ExecStart=/opt/fram/bin/fram-daemon ${FRAM_PORT} ${FRAM_LOG}
Restart=on-failure
RestartSec=2

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/north

[Install]
WantedBy=multi-user.target
