# deploy/tern-gateway.service
# systemd unit — the authenticated multi-tenant gateway.
#   sudo cp deploy/tern-gateway.service /etc/systemd/system/
#   printf 'GATEWAY_PORT=8088\nGATEWAY_TENANTS=/var/lib/tern/tenants.edn\n' \
#     | sudo tee /etc/tern/gateway.env
#   sudo systemctl enable --now tern-gateway
# Then terminate TLS in a reverse proxy (Caddy/nginx) in front of GATEWAY_PORT.
[Unit]
Description=Tern auth gateway
After=network.target

[Service]
Type=simple
User=tern
Group=tern
Environment=FRAM_HOME=/opt/fram
EnvironmentFile=/etc/tern/gateway.env
# adjust if babashka is elsewhere (`command -v bb`)
ExecStart=/usr/local/bin/bb -cp /opt/tern/out /opt/tern/deploy/gateway/gateway.clj
Restart=on-failure
RestartSec=2

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/opt/tern
ReadWritePaths=/var/lib/tern

[Install]
WantedBy=multi-user.target
